# setup-guide/skill-templates/consolidating/consolidating.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Consolidating (weekly memory curator) launchd job template — WATNEY / Personal AI Kit

  Weekly weigh-in that flags when long-term memory is bloating. Invokes the
  consolidating skill via `claude -p` — same robust pattern as the dreaming job.
  This means the skill runs even if scripts/curate.py hasn't been built yet
  (the AI does the curation reasoning directly), and it stays consistent with
  the dreaming job so the canonical-path fix applies to both uniformly.

  Replace [USER] with the macOS username and [AI_NAME] with the lowercase AI
  name before installing. Fires weekly: Sunday 10:30.

  Install (setup.sh does this automatically):
    1. Save as: ~/Library/LaunchAgents/com.[user].[ai-name].consolidating.plist
    2. Run: launchctl load ~/Library/LaunchAgents/com.[user].[ai-name].consolidating.plist
    3. Verify: launchctl list | grep consolidating   (exit code should be 0)

  NOTE: the `cd "$HOME/[AI_NAME]"` path here matches the dreaming job. If the
  kit's canonical install path changes (pending TCC/launchd fix), update BOTH
  this and dreaming.plist.template together.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.[USER].[AI_NAME].consolidating</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$HOME/[AI_NAME]" &amp;&amp; /opt/homebrew/bin/claude -p "Run the consolidating skill — the weekly memory curator. Check whether long-term memory has crept past a healthy size, flag stale or duplicate entries, and report what needs a trim. REPORT ONLY — never rewrite memory without approval." 2&gt;&gt; "$HOME/[AI_NAME]/logs/consolidating.err"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/[USER]</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/[USER]/[AI_NAME]/logs/consolidating.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/[USER]/[AI_NAME]/logs/consolidating.err</string>
</dict>
</plist>
